# rootfs/etc/services.d/cron/run
#!/usr/bin/env sh

sudo -E -u azuracast jobberrunner -p 9004 /var/azuracast/.jobber