# rootfs/etc/services.d/php-fpm/run
#!/bin/sh

/usr/sbin/php-fpm7 -F --fpm-config /etc/php7/php-fpm.d/www.conf -c /etc/php7/